# images/keepalived/runit/etc/rc.local
#!/bin/sh
# create a directory to put enabled service files
mkdir /etc/service/enabled

cp -a /etc/service/available/confd /etc/service/enabled/
cp -a /etc/service/available/keepalived /etc/service/enabled/
cp -a /etc/probe.sh  /etc/confd/templates/probe.sh

